# com.raihan.garmin-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.raihan.garmin-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/raihan/training/.venv/bin/python</string>
        <string>-m</string>
        <string>garmin</string>
        <string>sync</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/raihan/training</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/raihan/training/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/raihan/training/sync.log</string>
</dict>
</plist>
